# systemd/usbreset@.timer
[Unit]
Description=usbreset: Periodically reset usb device

[Timer]
OnBootSec=120min
OnUnitActiveSec=120min

[Install]
WantedBy=timers.target
