# systemd/usbreset@.timer
[Unit]
Description=usbreset: Periodically reset usb device

[Timer]
OnBootSec=240min
OnUnitActiveSec=240min

[Install]
WantedBy=timers.target
